# vm-backup.timer
[Unit]
Description=VM backup every day at 0200

[Timer]
OnCalendar=*-*-* 02:00:00
Unit=vm-backup.service

[Install]
WantedBy=multi-user.target
